# image/php/supervisord.conf
[supervisord]
nodaemon=true

[program:php_fpm]
command=php-fpm

[program:laravel-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan queue:work --sleep=3 --tries=2
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/supervisor_queue-work.log